# contrib/ym-bridge.service
[Unit]
Description=Yandex Music MPRIS Bridge
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ym-bridge run
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
